# doc/examples/crontab
0 5 * * 1 /usr/local/bin/build-system-builder
